# pythonapp.conf
[program:pythonapp]
command=/usr/bin/uwsgi --ini /etc/uwsgi/apps-enabled/pythonapp.ini --uid python --gid python
autorestart=true
stopsignal=INT
